# deploy/telegram-max-bridge-web.service
[Unit]
Description=Telegram -> MAX bridge web UI
After=network-online.target

[Service]
Type=simple
User=bridge
WorkingDirectory=/opt/telegram-max-bridge
ExecStart=/opt/telegram-max-bridge/.venv/bin/python -m bridge.webui
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
